# supervisor/supervisord.conf
[supervisord]
logfile=/dev/null
pidfile=/tmp/supervisord.pid
nodaemon=true
strip_ansi=false
childlogdir=/tmp

[program:syncd]
command=/bin/bash -lc '/home/user/scripts/run-syncd.sh'
directory=/home/user
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=1
priority=10
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
directory=/CLIProxyAPI
autostart=true
autorestart=true
priority=20
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:cli-proxy-api]
command=/bin/bash -lc '/home/user/scripts/wait-for-sync.sh && /home/user/scripts/run-cliproxyapi.sh'
directory=/CLIProxyAPI
autostart=true
autorestart=true
priority=30
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:filebrowser]
command=/bin/bash -lc '/home/user/scripts/wait-for-sync.sh && /home/user/filebrowser --address 0.0.0.0 --port 8888 --root / --baseurl /filebrowser --database /home/user/filebrowser-data/filebrowser.db'
directory=/home/user
autostart=true
autorestart=true
priority=40
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:gotty]
command=/bin/bash -lc '/home/user/gotty --address 127.0.0.1 --port 8080 --path /t --credential ${GOTTY_USERNAME:-admin}:${GOTTY_PASSWORD:-adminadminadmin} -w bash'
directory=/home/user
autostart=true
autorestart=true
priority=45
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
