# supervisor/supervisord.conf
[supervisord]
logfile=/dev/null
pidfile=/home/user/supervisord.pid
nodaemon=true
strip_ansi=false
childlogdir=/home/user/logs

[program:syncd]
command=/bin/bash -lc '/home/user/scripts/run-syncd.sh'
directory=/home/user
autostart=true
autorestart=unexpected
startsecs=1
exitcodes=0
priority=10
user=root
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:cli-proxy-api]
command=/bin/bash -lc '/home/user/scripts/wait-for-sync.sh && /home/user/scripts/run-cliproxyapi.sh'
directory=/CLIProxyAPI
autostart=true
autorestart=true
user=root
priority=20
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:cpa-manager-plus]
command=/bin/bash -lc '/home/user/scripts/wait-for-sync.sh && /home/user/scripts/run-cpa-manager-plus.sh'
directory=/home/user
autostart=true
autorestart=true
user=root
priority=25
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:cpa-usage-keeper]
command=/bin/bash -lc '/home/user/scripts/wait-for-sync.sh && /home/user/scripts/run-cpa-usage-keeper.sh'
directory=/home/user
autostart=true
autorestart=unexpected
exitcodes=0
user=root
priority=27
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:filebrowser]
command=/bin/bash -lc '/home/user/scripts/wait-for-sync.sh && /home/user/filebrowser --address 0.0.0.0 --port 8888 --root / --database /home/user/filebrowser-data/filebrowser.db'
directory=/home/user
autostart=true
autorestart=true
priority=30
user=root
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:gotty]
command=/bin/bash -lc '/home/user/gotty --address 0.0.0.0 --port 18080 --credential ${GOTTY_USERNAME:-admin}:${GOTTY_PASSWORD:-adminadminadmin} -w /bin/bash'
directory=/home/user
autostart=true
user=root
autorestart=true
priority=40
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
